# haproxy/supervisord.conf
[supervisord]
nodaemon=true

[program:consultemplate]
command=consul-template -config=/etc/haproxy.json -consul=consul:8500