# haproxy/supervisord.conf
[supervisord]
nodaemon=true

[program:haproxy]
command=haproxy -f /etc/haproxy/haproxy.cfg

[program:consultemplate]
command=consul-template -config=/etc/haproxy.json -consul=192.168.99.100:8500